# deploy/supervisord.conf
; supervisord runs as PID 1 and manages the three long-lived processes.
; Program stdout/stderr are redirected to PID 1's fds so all logs surface in
; `docker logs` / the Dokploy log view.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; --- the Telegram bot (runs from the image at /app) ---
[program:bot]
command=npm start
directory=/app
autostart=true
autorestart=true
startretries=999
stopwaitsecs=10
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0

; --- headless Obsidian Sync (mirrors the host's obsidian-sync.service) ---
; Until `ob login` + `ob sync-setup` have been run once on the mounted config
; volume, this will fail and retry; that's expected and harmless.
[program:obsidian-sync]
command=ob sync --continuous
directory=/root/vaults
autostart=true
autorestart=true
startretries=999
startsecs=5
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0

; --- zombie-poller watchdog: restarts the bot if it stops draining updates ---
[program:watchdog]
command=/app/deploy/watchdog.sh
autostart=true
autorestart=true
startretries=999
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0
